# DeviceBoxApp/DeviceBox.service
[Unit]
Description=Device Box

[Service]
ExecStart=%h/bin/DeviceBoxApp localhost
Restart=always

[Install]
WantedBy=default.target
